# roles/vnc/templates/vncserver@.service.j2
[Unit]
Description=Start TigerVNC server at startup for display :%i
After=network.target
Requires=systemd-user-sessions.service

[Service]
Type=forking
User={{ vnc_user }}
Group={{ vnc_user }}
WorkingDirectory=/home/{{ vnc_user }}

PAMName=login
Environment="HOME=/home/{{ vnc_user }}"
Environment="DISPLAY=:%i"
Environment="XDG_SESSION_TYPE=x11"
Environment="XDG_RUNTIME_DIR=/run/user/{{ vnc_uid | default('1000') }}"
Environment="XAUTHORITY=/home/{{ vnc_user }}/.Xauthority"

ExecStartPre=/usr/bin/vncserver -kill :%i > /dev/null 2>&1 || true
ExecStart=/usr/bin/vncserver :%i -geometry 1920x1200 -dpi 96
ExecStop=/usr/bin/vncserver -kill :%i

Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
