# roles/vnc/templates/vncserver@.service.j2
[Unit]
Description=Start TigerVNC server at startup for display :%i
After=network.target
Requires=systemd-user-sessions.service

[Service]
Type=forking
User=dave
Group=dave
WorkingDirectory=/home/dave

PAMName=login
Environment="HOME=/home/dave"
Environment="DISPLAY=:%i"
Environment="XDG_SESSION_TYPE=x11"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="XAUTHORITY=/home/dave/.Xauthority"

ExecStartPre=/usr/bin/vncserver -kill :%i > /dev/null 2>&1 || true
ExecStart=/usr/bin/vncserver :%i -geometry 1920x1200 -dpi 96
ExecStop=/usr/bin/vncserver -kill :%i

Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
