# roles/fail2counter/templates/fail2counter.service.j2
[Unit]
Description=Fail2Counter IP scanner with AI exploit analysis
After=network.target webdis.service postgresql.service fail2counter-vpn.service
Requires=fail2counter-vpn.service

[Service]
ExecStart=/opt/fail2counter/fail2counter_worker.py
Restart=always
RestartSec=10
User=postgres
Group=postgres

# Redis credentials
Environment=REDIS_PASSWORD={{ redis_password }}

# PostgreSQL connection (uses MSF's Postgres instance via peer auth)
Environment="FAIL2COUNTER_DSN=host=/var/run/postgresql dbname=fail2counter user=postgres"

# Google Vertex AI / Claude configuration
Environment=ANTHROPIC_VERTEX_PROJECT_ID={{ anthropic_vertex_project_id | default('itpc-gcp-ansible-pe-eng-claude') }}
Environment=ANTHROPIC_VERTEX_REGION={{ anthropic_vertex_region | default('us-east5') }}
Environment=GOOGLE_APPLICATION_CREDENTIALS={{ google_application_credentials | default('/opt/fail2counter/gcp-sa-key.json') }}
{% if anthropic_vertex_model is defined %}
Environment=ANTHROPIC_VERTEX_MODEL={{ anthropic_vertex_model }}
{% endif %}

# Allow network namespace access
AmbientCapabilities=CAP_SYS_ADMIN CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
