# roles/rclone/files/rclone-backup.service
[Unit]
Description=Rclone backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rclone-backup.sh
